# src/systemd/tfgrid-ai-manager@.service
[Unit]
Description=TFGrid AI Agent Manager Handler (connection %i)

[Service]
Type=simple
ExecStart=/opt/ai-agent/scripts/socket-handler.sh
StandardInput=socket
StandardOutput=socket
StandardError=journal
KillMode=process
User=root
MemoryMax=50M
CPUQuota=5%
