# src/systemd/tfgrid-ai-manager.socket
[Unit]
Description=TFGrid AI Agent Manager Socket

[Socket]
ListenStream=/run/ai-agent.sock
SocketMode=0600
Accept=yes

[Install]
WantedBy=sockets.target
